# supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=mysqld
autorestart=true
priority=1

[program:cron]
command=cron
autorestart=true
priority=10
